# 7.0/services/php-fpm/run
#!/bin/sh

# run service
exec /usr/local/sbin/php-fpm -c /etc/php/ -y /etc/php/php-fpm.conf
